# lib/startup.conf
[program:startup]
command=/usr/local/bin/startup.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/startup.err.log
stdout_logfile=/var/log/startup.out.log
